# playbooks/templates/minecraft.service
[Unit]
Description=The Minecraft Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/mnt/volume-tassadar/minecraft/
ExecStart=java -Xmx1024M -Xms512M -jar server.jar nogui
Restart=on-failure

[Install]
WantedBy=multi-user.target